# buildroot-external/overlay/base/etc/profile.d/ulimits.sh
#!/bin/sh
# shellcheck shell=dash
ulimit -S -c unlimited >/dev/null 2>&1
