# buildroot-external/overlay/base/etc/profile.d/umask.sh
#!/bin/sh
# shellcheck shell=dash
umask 0002
