# buildroot-external/package/recovery-system/external/overlay/etc/profile.d/aliases.sh
alias ll='ls -alF'
alias la='ls -A'
alias l='ls -CF'
