# buildroot-external/overlay/base-raspmatic_oci/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Init the system
null::sysinit:/bin/mkdir -p /dev_host
#Bug in Kubernetes -> sys is mounted RO even for privileged -> https://github.com/kubernetes/kubernetes/pull/96877
null::sysinit:/bin/mount -o rw,remount /sys
null::sysinit:/bin/mount -t tmpfs tmpfs /tmp

# run any rc scripts (start) and move boot.log
::sysinit:/etc/init.d/rcS 2>&1 | /usr/bin/tee -a /tmp/boot.log
::sysinit:/bin/mv /tmp/boot.log /var/log/boot.log

# run monit to monitor our important services
null::respawn:/usr/bin/monit -Ic /etc/monitrc

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
null::shutdown:/usr/bin/monit unmonitor all

# run any rc scripts (stop)
::shutdown:/etc/init.d/rcK
